# deploy/systemd/trace.service
[Unit]
Description=Trace 集装箱轨迹批处理
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/trace
Environment=TRACE_PROJECT_DIR=/opt/trace
EnvironmentFile=-/etc/sysconfig/trace
ExecStart=/opt/trace/deploy/run-trace.sh
# CentOS 7 systemd uses 0 to disable the start timeout.
TimeoutStartSec=0
Nice=10
